# deploy/azure-ubuntu/chatsense-stock-sync.service
[Unit]
Description=ChatSense stock data sync
After=network.target postgresql.service

[Service]
Type=oneshot
User=chatsense
Group=www-data
WorkingDirectory=/opt/chatsense/app
EnvironmentFile=/etc/chatsense/backend.env
ExecStart=/opt/chatsense/venv/bin/python manage.py sync_stock_data
